# fab_deploy2/default-configs/templates/base/gunicorn/upstart_gunicorn.conf
description "gunicorn server"
author "RED"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 10 5

script
sh {{ gunicorn.start_script }}
end script
